# workflows/bootc-debs-builder/openipmipostinstall/usr/lib/systemd/system/openipmi.service.d/skip-without-ipmi.conf
# Skip the OpenIPMI driver init when the platform declares no IPMI/BMC at all —
# i.e. no DMI type 38 (IPMI Device Information) entry, as on a virtual machine.
# The service is then cleanly "condition-skipped" instead of failing. If an IPMI
# device IS declared in DMI but the driver still fails to load, the condition
# passes, the service runs, and the real error is reported.
[Unit]
ConditionPathExistsGlob=/sys/firmware/dmi/entries/38-*
