# workflows/bootc-debs-builder/chronypostinstall/usr/lib/systemd/system/chrony.service.d/wait-network.conf
# chrony needs the network before it can reach any NTP server. The stock unit
# does not wait for it, so on a fast boot chronyd starts before networking is up
# and fails its initial sync. Order it after the network is online.
[Unit]
After=network-online.target
Wants=network-online.target
