# workflows/bootc-debs-builder/fanctl/usr/lib/systemd/system/fanctl.service
[Unit]
Description=fanctl — generic fan controller (local IPMI)
Documentation=https://github.com/DaemonCores/DaemonCores-VE
# Only start when the local IPMI interface is present.
ConditionPathExists=/dev/ipmi0
After=multi-user.target

[Service]
Type=simple
ExecStart=/usr/sbin/fanctl --daemon
# Safety net if the process is killed (SIGKILL cannot be trapped): explicitly
# hand fan regulation back to the BMC automatic mode.
ExecStopPost=/usr/sbin/fanctl --restore
Restart=on-failure
RestartSec=10
Nice=5

# Persistent calibration state in /var/lib/fanctl (created and made writable).
StateDirectory=fanctl

# Light hardening (needs /dev/ipmi0 and to execute ipmitool).
ProtectHome=yes
ProtectSystem=strict
NoNewPrivileges=yes

[Install]
WantedBy=multi-user.target
